# conf.d/x11vnc.conf
[program:x11vnc]
user=%(ENV_USER)s
command=x11vnc -display :123 -forever -shared -usepw
autorestart=true
