# conf.d/xfce4.conf
[program:xfce4]
user=%(ENV_USER)s
environment=
    DISPLAY=:123
command=xfce4-session
autorestart=true